# client-agent/systemd/m5nas-monitor.service
[Unit]
Description=M5Stack NAS Monitor Daemon
Documentation=https://github.com/example/m5stack-nas-monitor
After=network.target zfs.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/m5nas-monitor --foreground --config /etc/m5nas-monitor.conf
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=m5nas-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run

# Resource limits
LimitNOFILE=1024
MemoryMax=64M

[Install]
WantedBy=multi-user.target